# usb_disconnect.service
[Unit]
Description=USB Disconnect Handler

[Service]
Type=oneshot
ExecStart=/scripts/usb_disconnect.sh

[Install]
WantedBy=multi-user.target

